# systemd/radarr.service
[Unit]
Description=Radarr MVP - Movie Management System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/radarr
ExecStart=/opt/radarr/radarr-mvp
EnvironmentFile=/opt/radarr/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radarr

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target